# etc/supervisord.d/fhgfs-admon.ini
[program:fhgfs-admon]
command = /opt/fhgfs/sbin/fhgfs-admon cfgFile=/etc/fhgfs/fhgfs-admon.conf pidFile=/var/run/fhgfs-admon.pid
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
